# unfinished/systemd-etherpad.service
[Unit]
Description=etherpad service
After=network.target

[Service]
Type=simple
User=etherpad
ExecStart=/home/etherpad/etherpad/bin/run.sh
Restart=on-abort

